# config/chroot_local-includes/lib/systemd/system/NetworkManager.service.d/hardening.conf
[Service]

# Paths
ProtectProc=invisible

# Security
NoNewPrivileges=true

# Sandboxing
ProtectSystem=strict
ProtectHome=true
RuntimeDirectory=NetworkManager
StateDirectory=NetworkManager
ConfigurationDirectory=NetworkManager
PrivateTmp=true
PrivateDevices=true
ProtectHostname=true
ProtectClock=true
ProtectKernelLogs=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true

# System Call Filtering
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
