# config/chroot_local-includes/usr/lib/systemd/user/tails-start-system-gnome-session-target.service
[Unit]
Description=Start the system GNOME session target
After=gnome-session.target

[Service]
ExecStart=/bin/sudo /bin/systemctl start gnome-session.target

[Install]
WantedBy=gnome-session.target
